# assets/launchd/com.remote.bastion-session.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.remote.bastion-session</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/bastion-session</string>
      <string>watch</string>
      <string>--interval</string>
      <string>300</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PROFILE</key>
      <string>oabcs1-terraform</string>
      <key>REGION</key>
      <string>us-chicago-1</string>
      <key>AUTH_METHOD</key>
      <string>security_token</string>
      <key>TARGET_USER</key>
      <string>opc</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/bastion-session.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/bastion-session.err.log</string>
  </dict>
</plist>
